# ftpup.service
[Unit]
Description=ftpup

[Service]
ExecStart=/usr/local/bin/ftpup

[Install]
WantedBy=multi-user.target
